# profile.d/mysql.sh
export MYSQL_HOME=/usr/local/lnmp/mysql
export PATH=$MYSQL_HOME/bin:$PATH
