# profile.d/php.sh
export PHP_HOME=/usr/local/lnmp/php56
export PATH=$PHP_HOME/bin:$PATH
