# profile.d/redis.sh
export REDIS_HOME=/usr/local/lnmp/redis
export PATH=$REDIS_HOME/bin:$PATH
